# root/etc/services.d/timemachine/run
#!/usr/bin/with-contenv bash

echo "Starting Time Machine."
exec \
    s6-setuidgid abc /bin/bash -c \
    '/usr/sbin/netatalk -d -F /config/afp.conf'