# root/etc/services.d/timemachine/run
#!/usr/bin/with-contenv bash

echo "Starting Time Machine."
exec \
    /usr/sbin/netatalk -d -F /config/afp.conf